# god/templates/default/god.initrc.erb
#!/bin/sh
#
# God
#
# chkconfig: - 85 15
# description: start, stop, restart God
#

GODDIR=`gem env | grep "EXECUTABLE DIRECTORY"| cut -d " " -f 6`
GODCMD=${GODDIR}/god
PIDFILE=/var/run/god.pid
GODCONF=/usr/local/etc/god.conf
name="god"

#FreeBSD compatibility block
#don't forget  about rc.conf (god_enable="YES")
if [ -f /etc/rc.subr ]; then
    . /etc/rc.subr
    rcvar=`set_rcvar`
fi
 
RETVAL=0
 
case "$1" in
    start)
      $GODCMD -P $PIDFILE -c $GODCONF
      RETVAL=$?
  ;;
    stop)
      $GODCMD terminate || if [ -f $PIDFILE ]; then 
				echo "Sending KILL signal..."
				kill -9 `cat ${PIDFILE}`
				rm $PIDFILE
			   fi
      RETVAL=$?
  ;;
    restart)
      $GODCMD terminate || kill -9 `cat ${PIDFILE}`
      $GODCMD -P $PIDFILE -c $GODCONF
      RETVAL=$?
  ;;
    status)
      $GODCMD status && echo "[ OK ] god is up and running"
      RETVAL=$?
  ;;
    *)
      echo "Usage: god {start|stop|restart|status}"
      exit 1
  ;;
esac
 
exit $RETVAL
